# restore_host_ssh.service
[Unit]
Description=One time boot script
[Service]
Type=simple
ExecStart=/restore_host_ssh.sh
[Install]
WantedBy=multi-user.target 
